# server-code/vr-start.service
[Unit]
Description=Vena's Radio Service
After=network.target

[Service]
User=vr
WorkingDirectory=/root/Venas-Radio
ExecStart=/usr/bin/env bash -c "source /root/.nvm/nvm.sh && yarn run start"
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target